# setup/jarvis-provisioning.service
[Unit]
Description=Jarvis Node Provisioning Server
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKING_DIR__
Environment=CONFIG_PATH=__WORKING_DIR__/config.json
Environment=JARVIS_SECRET_DIRECTORY=__SECRET_DIR__
ExecStart=/usr/bin/python3 scripts/run_provisioning.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
